# realty-main/deploy/etc/s6-overlay/s6-rc.d/workers/run
#!/command/with-contenv execlineb

cd /app
exec python -m realty db_worker -v 3
