# realty-main/deploy/etc/s6-overlay/s6-rc.d/server/run
#!/command/with-contenv execlineb

cd /app
exec python -m realty prodserver --workers 4
